# pkgbuild/mhwd-biglinux.install
post_install() {
    systemctl enable nvidia-manager-verify
}

pre_remove() {
    systemctl disable nvidia-manager-verify
}
